# scripts/com.synapse.open_daily_urls.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.synapse.open_daily_urls</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/lei/dev/personal/synapse/scripts/open_daily_urls.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/lei/dev/personal/synapse/scripts/open_daily_urls.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/lei/dev/personal/synapse/scripts/open_daily_urls.err</string>
</dict>
</plist>
